# etc/systemd/system/omv-dropbox-backup.service
[Unit]
Description=OpenMediaVault Dropbox Backup
Documentation=https://github.com/blortorbis/openmediavault-dropbox-backup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/omv-dropbox-backup
StandardOutput=journal
StandardError=journal

# Security hardening.
# Note: ProtectSystem=full (not strict) because rclone must be able to
# write to the backup destination under /srv, and ProtectHome must stay
# off because rclone refreshes its OAuth token in /root/.config/rclone.
NoNewPrivileges=yes
ProtectSystem=full
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
