# etc/systemd/system/omv-dropbox-backup.timer
[Unit]
Description=OpenMediaVault Dropbox Backup Timer
Documentation=https://github.com/blortorbis/openmediavault-dropbox-backup

[Timer]
# The schedule is overridden by a drop-in generated from the plugin
# configuration (see the Salt state). Default: daily at 2:00 AM.
OnCalendar=*-*-* 02:00:00
RandomizedDelaySec=1800
Persistent=true

[Install]
WantedBy=timers.target
